# .autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crypto-campaign.auto-deploy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/Danallovertheplace/crypto-campaign-unified/scripts/auto-start-on-directory-open.js</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/Danallovertheplace/crypto-campaign-unified/.autostart-system.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/Danallovertheplace/crypto-campaign-unified/.autostart-system.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/Danallovertheplace/crypto-campaign-unified</string>
    
    <key>WatchPaths</key>
    <array>
        <string>/Users/Danallovertheplace/crypto-campaign-unified/frontend/src</string>
        <string>/Users/Danallovertheplace/crypto-campaign-unified/backend/src</string>
        <string>/Users/Danallovertheplace/crypto-campaign-unified/contracts/src</string>
    </array>
    
    <key>QueueDirectories</key>
    <array>
        <string>/Users/Danallovertheplace/crypto-campaign-unified</string>
    </array>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>